# src/etc/systemd/system/bluezalsa.service
[Unit]
Description=BluezAlsa proxy
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bluealsa --profile=a2dp-source
Restart=on-failure
RestartSec=5

# Sandboxing — conservative because bluealsa runs as root and needs
# direct Bluetooth hardware access via AF_BLUETOOTH kernel sockets.
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
